# meta-citadel/recipes-citadel/citadel-config/files/systemd/launch-default-realm.path
[Unit]
Description=Wait for creation of PA and Wayland sockets before launching first realm.

[Path]
PathChanged=/run/user/1000/pulse/native
PathChanged=/run/user/1000/wayland-0
